# exam/roles/influxdb/templates/testPing.service.j2
[Unit]
Description=testPing

[Service]
User=user
ExecStart=/opt/testPing.sh

[Install]
WantedBy=multi-user.target
