# exam/roles/apache_exporter/templates/apache_exporter.service.j2
[Unit]
Description=Apache Exporter

[Service]
User=root
Restart=always
ExecStart=/usr/local/bin/apache_exporter \

[Install]
WantedBy=multi-user.target
